# build-tools/gauffr.cron.daily
#!/bin/bash
#
# Gauffr's cron script
#

/usr/bin/gauffr_clear_log

exit 0